# DEBIAN/prerm
#!/bin/bash

set -e

systemctl disable --now delete-old-downloads.timer || true
systemctl disable --now delete-old-downloads.service || true

exit 0
