# etc/systemd/system/delete-old-downloads.timer
[Unit]
Description=Run Clean Downloads script hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
